# systemd/audiobook-proxy.service
[Unit]
Description=Audiobooks Library HTTPS Reverse Proxy
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
After=network-online.target audiobook-api.service
Requires=audiobook-api.service
Wants=network-online.target
# Ensure /opt BTRFS subvolume is mounted before starting (prevents intermittent CHDIR failures)
RequiresMountsFor=/opt/audiobooks

[Service]
Type=simple
User=audiobooks
Group=audiobooks
WorkingDirectory=/opt/audiobooks/library/web-v2

# Load configuration from /etc/audiobooks/audiobooks.conf
EnvironmentFile=/etc/audiobooks/audiobooks.conf
Environment=PYTHONUNBUFFERED=1

# Only start if port is not already in use
ExecStartPre=/bin/sh -c '! /usr/bin/ss -tlnp | /usr/bin/grep -q ":${AUDIOBOOKS_WEB_PORT} "'
ExecStart=/usr/bin/python3 /opt/audiobooks/library/web-v2/proxy_server.py

Restart=on-failure
RestartSec=5
# Use journal logging (query with: journalctl -u audiobooks-proxy)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobooks-proxy

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/audiobooks /opt/audiobooks
# Data directory for audiobooks database and cache
ReadWritePaths=/var/lib/audiobooks
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=audiobook.target multi-user.target
