# systemd/audiobook.target
[Unit]
Description=Audiobook Library Services
Documentation=https://github.com/greogory/Audiobook-Manager
After=network-online.target
Wants=network-online.target

# Core services - these are the main services that must start
Wants=audiobook-api.service
Wants=audiobook-proxy.service
Wants=audiobook-redirect.service

# Background workers - conversion and file moving
Wants=audiobook-converter.service
Wants=audiobook-mover.service

# Scheduled tasks - timers trigger their corresponding .service units
Wants=audiobook-downloader.timer

[Install]
WantedBy=multi-user.target
